# autostart/feh.desktop
[Desktop Entry]
Type=Aplication
Name=feh
Exec=/home/pi/Desktop/start.sh
